# mongo/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:mongodb]
command=mongod
autostart=true
autorestart=true
